# src/crontab
# START CRON JOB
0 8 * * 1-5 /usr/local/bin/python /app/App.py --config /app/config/config.ini --filter /app/config/filters.conf >> /log/app.log 2>&1
# END CRON JOB